# server-configs/aimaid.service
[Unit]
Description=Aimaid data and WebRTC signaling service
After=network-online.target mysqld.service aiprovider-ccxt.service
Wants=network-online.target
Requires=aiprovider-ccxt.service

[Service]
Type=simple
User=aimaid
Group=aimaid
WorkingDirectory=/opt/aimaid
EnvironmentFile=/etc/aimaid/aimaid.env
ExecStart=/usr/bin/java -XX:MaxRAMPercentage=55 -jar /opt/aimaid/ai-provider.jar
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aimaid

[Install]
WantedBy=multi-user.target
