# server-configs/aiprovider-ccxt.service
[Unit]
Description=Aimaid CCXT read-only market data gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aimaid
Group=aimaid
WorkingDirectory=/opt/aimaid/ccxt-gateway
EnvironmentFile=/etc/aimaid/aimaid.env
ExecStart=/usr/bin/node /opt/aimaid/ccxt-gateway/src/server.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
